# debian/postinst
#!/bin/bash
# postinst script for stumpfworks-nas
# This script runs after package installation

set -e

case "$1" in
    configure)
        # Create data directory
        mkdir -p /var/lib/stumpfworks-nas
        chmod 755 /var/lib/stumpfworks-nas

        # Create config directory
        mkdir -p /etc/stumpfworks-nas
        chmod 755 /etc/stumpfworks-nas

        # Initialize database if not exists
        if [ ! -f /var/lib/stumpfworks-nas/stumpfworks.db ]; then
            echo "Database will be initialized on first start..."
        fi

        # Reload systemd
        systemctl daemon-reload

        # Enable service (don't start automatically - let admin configure first)
        systemctl enable stumpfworks-nas.service

        # Print post-install message
        echo ""
        echo "╔══════════════════════════════════════════════════════════════╗"
        echo "║         Stumpf.Works NAS v1.0.0 installed successfully!      ║"
        echo "╚══════════════════════════════════════════════════════════════╝"
        echo ""
        echo "Next steps:"
        echo ""
        echo "1. Configure the system:"
        echo "   sudo nano /etc/stumpfworks-nas/config.yaml"
        echo ""
        echo "2. Start the service:"
        echo "   sudo systemctl start stumpfworks-nas"
        echo ""
        echo "3. Check status:"
        echo "   sudo systemctl status stumpfworks-nas"
        echo ""
        echo "4. Access the web interface:"
        echo "   http://YOUR_SERVER_IP:8080"
        echo ""
        echo "5. Default credentials:"
        echo "   Username: admin"
        echo "   Password: admin"
        echo "   (Change immediately after first login!)"
        echo ""
        echo "For more information:"
        echo "  - Documentation: /usr/share/doc/stumpfworks-nas/"
        echo "  - Logs: sudo journalctl -u stumpfworks-nas -f"
        echo "  - GitHub: https://github.com/Stumpf-works/stumpfworks-nas"
        echo ""
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
