# debian/postinst
#!/bin/bash
# postinst script for stumpfworks-nas
# This script runs after package installation

set -e

case "$1" in
    configure)
        echo "â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”"
        echo "  StumpfWorks NAS - Post-Installation Setup"
        echo "â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”â”"
        echo ""

        # Create directories
        echo "ðŸ“ Creating directories..."
        mkdir -p /var/lib/stumpfworks-nas/backups
        mkdir -p /var/lib/stumpfworks-nas/plugins
        mkdir -p /var/log/stumpfworks-nas
        mkdir -p /etc/stumpfworks-nas
        mkdir -p /usr/share/stumpfworks-nas/doc

        chmod 755 /var/lib/stumpfworks-nas
        chmod 755 /etc/stumpfworks-nas
        chmod 755 /var/log/stumpfworks-nas

        # Setup PostgreSQL
        echo "ðŸ—„ï¸  Setting up PostgreSQL database..."

        # Check if PostgreSQL is running
        if ! systemctl is-active --quiet postgresql; then
            echo "   Starting PostgreSQL..."
            systemctl start postgresql
        fi

        # Run database setup as postgres user
        if ! sudo -u postgres psql -lqt | cut -d \| -f 1 | grep -qw stumpfworks_nas; then
            echo "   Creating database and user..."

            # Generate secure password
            DB_PASSWORD=$(openssl rand -base64 32 | tr -d "=+/" | cut -c1-32)

            # Create password file (must be done as root before running as postgres user)
            echo "$DB_PASSWORD" > /etc/stumpfworks-nas/.db-password
            chmod 600 /etc/stumpfworks-nas/.db-password

            # Create database and user
            sudo -u postgres psql <<EOF
CREATE USER stumpfworks WITH PASSWORD '$DB_PASSWORD';
CREATE DATABASE stumpfworks_nas OWNER stumpfworks;
GRANT ALL PRIVILEGES ON DATABASE stumpfworks_nas TO stumpfworks;
EOF

            echo "   âœ“ Database created"
        else
            echo "   âœ“ Database already exists"
        fi

        # Create config if it doesn't exist
        if [ ! -f /etc/stumpfworks-nas/config.yaml ]; then
            echo "ðŸ“ Creating default configuration..."

            # Read DB password
            if [ -f /etc/stumpfworks-nas/.db-password ]; then
                DB_PASSWORD=$(cat /etc/stumpfworks-nas/.db-password)
            else
                DB_PASSWORD=""
            fi

            # Generate JWT secret
            JWT_SECRET=$(openssl rand -base64 32)

            cat > /etc/stumpfworks-nas/config.yaml <<CONFIGEOF
# Stumpf.Works NAS Configuration
app:
  name: "Stumpf.Works NAS"
  version: "0.1.0"
  environment: "production"

server:
  host: "0.0.0.0"
  port: 8080
  readTimeout: "15s"
  writeTimeout: "15s"
  idleTimeout: "60s"

database:
  driver: "postgres"
  host: "localhost"
  port: 5432
  database: "stumpfworks_nas"
  username: "stumpfworks"
  password: "$DB_PASSWORD"
  sslmode: "disable"
  maxOpenConns: 25
  maxIdleConns: 5
  connMaxLifetime: "5m"

auth:
  jwtSecret: "$JWT_SECRET"
  jwtExpirationHours: 24
  jwtRefreshHours: 168
  bcryptCost: 12
  sessionTimeout: "24h"

logging:
  level: "info"
  development: false
CONFIGEOF

            chmod 600 /etc/stumpfworks-nas/config.yaml
            echo "   âœ“ Configuration created"
        else
            echo "   âš ï¸  Existing config found, keeping it"
        fi

        # Setup Samba for network shares
        echo "ðŸŒ Configuring Samba for Windows Network Discovery..."

        # Create shares directory
        mkdir -p /etc/samba/shares.d
        chmod 755 /etc/samba/shares.d

        # Only create/update Samba config if it doesn't have our marker or doesn't exist
        if ! grep -q "Stumpf.Works NAS" /etc/samba/smb.conf 2>/dev/null; then
            # Backup existing smb.conf if it exists
            if [ -f /etc/samba/smb.conf ]; then
                cp /etc/samba/smb.conf /etc/samba/smb.conf.backup-$(date +%Y%m%d-%H%M%S) 2>/dev/null || true
            fi

            echo "   Creating Samba configuration..."
            # Create Samba base configuration with Windows Network Discovery
            cat > /etc/samba/smb.conf <<'SMBCONF'
#======================= Global Settings =======================

[global]
   # Basic Server Settings
   workgroup = WORKGROUP
   server string = Stumpf.Works NAS
   netbios name = STUMPFWORKS-NAS
   server role = standalone server

   # Logging
   log file = /var/log/samba/log.%m
   max log size = 1000
   logging = file
   panic action = /usr/share/samba/panic-action %d

   # Windows Network Discovery / Browser Settings
   # This makes the NAS visible in Windows Network Explorer
   local master = yes
   preferred master = yes
   os level = 65
   domain master = no
   wins support = yes
   dns proxy = no

   # Name Resolution Order
   name resolve order = wins bcast host lmhosts

   # Enable SMB2/SMB3 for Windows 10/11 compatibility
   server min protocol = SMB2
   server max protocol = SMB3

   # Authentication
   security = user
   passdb backend = tdbsam
   obey pam restrictions = yes
   unix password sync = yes
   passwd program = /usr/bin/passwd %u
   passwd chat = *Enter\snew\s*\spassword:* %n\n *Retype\snew\s*\spassword:* %n\n *password\supdated\ssuccessfully* .
   pam password change = yes
   map to guest = bad user

   # Performance optimizations
   socket options = TCP_NODELAY IPTOS_LOWDELAY SO_RCVBUF=131072 SO_SNDBUF=131072
   read raw = yes
   write raw = yes
   max xmit = 65535
   min receivefile size = 16384
   use sendfile = yes
   aio read size = 16384
   aio write size = 16384

   # Character set
   unix charset = UTF-8
   dos charset = CP850

   # Include share configurations from shares.d
   include = /etc/samba/shares.d/*.conf

#======================= Share Definitions =======================
# Individual shares are configured in /etc/samba/shares.d/
# This allows the Stumpf.Works NAS backend to manage shares dynamically
SMBCONF
        else
            echo "   âœ“ Samba configuration already exists, keeping it"
        fi

        # Create SMB users group
        echo "ðŸ‘¥ Creating Samba user group..."
        if ! getent group smbusers > /dev/null 2>&1; then
            groupadd smbusers
            echo "   âœ“ Created smbusers group"
        else
            echo "   âœ“ smbusers group already exists"
        fi

        # Enable and start Samba services
        systemctl enable smbd nmbd 2>/dev/null || true
        systemctl restart smbd nmbd 2>/dev/null || true

        # Enable WSDD for Windows 10/11 Network Discovery
        systemctl enable wsdd 2>/dev/null || true
        systemctl start wsdd 2>/dev/null || true

        # Enable Avahi for macOS/Linux Discovery
        systemctl enable avahi-daemon 2>/dev/null || true
        systemctl start avahi-daemon 2>/dev/null || true

        echo "   âœ“ Samba configured with Windows Network Discovery"
        echo "   âœ“ Server will appear as: STUMPFWORKS-NAS"

        # Reload systemd
        echo "ðŸ”„ Reloading systemd..."
        systemctl daemon-reload

        # Enable service
        echo "âœ… Enabling service..."
        systemctl enable stumpfworks-nas.service

        # Print post-install message
        echo ""
        echo "â•”â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•—"
        echo "â•‘     Stumpf.Works NAS v0.1.0 installed successfully!          â•‘"
        echo "â•šâ•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•"
        echo ""
        echo "Next steps:"
        echo ""
        echo "1. Start the service:"
        echo "   sudo systemctl start stumpfworks-nas"
        echo "   or"
        echo "   stumpfctl service start"
        echo ""
        echo "2. Check status:"
        echo "   stumpfctl service status"
        echo ""
        echo "3. Access the web interface:"
        echo "   http://YOUR_SERVER_IP:8080"
        echo ""
        echo "4. View logs:"
        echo "   stumpfctl logs -f"
        echo ""
        echo "5. Manage users:"
        echo "   stumpfctl user list"
        echo ""
        echo "For more information:"
        echo "  - Documentation: /usr/share/doc/stumpfworks-nas/"
        echo "  - GitHub: https://github.com/Stumpf-works/stumpfworks-nas"
        echo ""
        echo "ðŸ’¡ Tip: Run 'stumpfctl --help' to see all available commands"
        echo ""
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
