# debian/prerm
#!/bin/bash
# prerm script for stumpfworks-nas
# This script runs before package removal

set -e

case "$1" in
    remove|upgrade|deconfigure)
        # Stop service before removal
        if systemctl is-active --quiet stumpfworks-nas.service; then
            echo "Stopping Stumpf.Works NAS service..."
            systemctl stop stumpfworks-nas.service || true
        fi

        # Disable service
        if systemctl is-enabled --quiet stumpfworks-nas.service; then
            echo "Disabling Stumpf.Works NAS service..."
            systemctl disable stumpfworks-nas.service || true
        fi
        ;;

    failed-upgrade)
        ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
